# deploy/systemd/iis-core-worker.service.tmpl
[Unit]
Description=IIS Core Engine Worker (ForgeFlow)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__IIS_USER__
Group=__IIS_USER__
WorkingDirectory=__IIS_APP_DIR__
EnvironmentFile=__IIS_APP_DIR__/.env
ExecStart=__IIS_VENV_BIN__/python -m app.workers.pipeline_worker
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
